# containers/backup/isp-kpn copy/services/pppoe-server/pppoe-server
#!/sbin/openrc-run

# Description of the service
description="PPPoE Server"

# Command to start the PPPoE server
command="/usr/sbin/pppoe-server"
command_args="-I isp_out.6 -q /usr/sbin/pppd -Q /usr/sbin/pppoe -C kpn -L 192.168.1.1 -R 192.168.1.10 -N 10"

# This service requires the network interfaces to be up
depend() {
    need network
}

# Start the service
start() {
    ebegin "Starting PPPoE Server"

    start-stop-daemon --start --exec $command -- $command_args
    eend $?
}

# Stop the service
stop() {
    ebegin "Stopping PPPoE Server"
    start-stop-daemon --stop --exec $command
    eend $?
}

# Check the status of the service
status() {
    status_of_proc "$command" "PPPoE Server"
}